# Docker/allinone/etc/services.d/ozwadmin/finish
#!/usr/bin/execlineb -S0
echo "Shutting Down Container after ozwdaemon exited"
s6-svscanctl -t /var/run/s6/services